# mrtg.cron
#!/bin/sh
#########
# Cron-job that can use many configs
#########

# Initial values:
CONFIGDIR="/etc/mrtg"
CONFIGS="mrtg.cfg"

# Get service config
[ -f /etc/sysconfig/mrtg ] && . /etc/sysconfig/mrtg

for i in ${CONFIGS}; do
	LANG=C /usr/bin/mrtg ${CONFIGDIR}/${i} --logging /var/log/mrtg/$i.log --lock-file /var/lib/mrtg/${i}_l --confcache-file=/var/lib/mrtg/${i}.ok
done
